# home/.config/autostart/silence-fan.desktop
[Desktop Entry]
Type=Application
Name=silence-fan
Path=/bin
Exec=silence-fan
StartupNotify=false
Terminal=false
Hidden=false

