# root/etc/systemd/system/resume-fan.service
[Unit]
Description=Fanafter ession using slock for user %i
After=sleep.target

[Service]
ExecStart=/bin/silence-fan

[Install]
WantedBy=sleep.target
